# etc/systemd/system/nvidia-clock-cap.service
[Unit]
Description=Cap NVIDIA boost clocks at ~1100 MHz for lower thermals
After=multi-user.target

[Service]
Type=oneshot
ExecStart=-/usr/bin/nvidia-smi --lock-gpu-clocks=0,1101

[Install]
WantedBy=multi-user.target
