# etc/systemd/system/tailscale-autoheal.service
[Unit]
Description=Restart tailscaled when node is offline or control map is stale
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c '\
  set -euo pipefail; \
  status_json="$(tailscale status --json 2>/dev/null || true)"; \
  if [ -z "$status_json" ]; then \
    logger -t tailscale-autoheal "tailscale status failed; restarting tailscaled"; \
    systemctl restart tailscaled; \
    exit 0; \
  fi; \
  if printf "%%s" "$status_json" | python3 -c "\
    import json, sys; \
    data = json.load(sys.stdin); \
    self_data = data.get(\"Self\", {}); \
    health = \" \".join(data.get(\"Health\", [])); \
    backend = data.get(\"BackendState\", \"\"); \
    offline = not bool(self_data.get(\"Online\")); \
    stale = \"network map\" in health.lower(); \
    bad_backend = backend not in (\"Running\", \"Starting\"); \
    raise SystemExit(0 if (offline or stale or bad_backend) else 1)"; then \
    logger -t tailscale-autoheal "tailscale unhealthy (offline/netmap/backend); restarting tailscaled"; \
    systemctl restart tailscaled; \
  fi'
